# server/build/fed/service/log/run
#!/bin/sh
LOG=/var/log/runsvdir/fed
test -d "$LOG" || mkdir -p -m755 "$LOG"
exec svlogd -t "$LOG"
